# dev-tools/ansible/roles/zabbix/files/check-needs-restarting.cron
#!/bin/sh

/usr/bin/needs-restarting -r > /dev/null
echo $? > /etc/zabbix/needs-restarting

exit 0
