# dev-tools/ansible/roles/tusd/templates/tus-temp-dir-cleanup.cron.j2
#!/bin/bash

# Delete files that are older than 7 days
find {{ tusd_upload_dir }} -mtime +7 -exec rm -rf {} \;
